# apv5sdk-v15/rootfs_s27_ramdisk/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K stop
#::shutdown:/etc/ath/killVAP all
::shutdown:/bin/umount -a
#::respawn:/sbin/getty ttyS0 115200
::respawn:/bin/ash