# apv5sdk-v15/autelan/base-files/files/usr/sbin/mld_snoop
#!/bin/sh /etc/rc.common

KILLALL=/usr/bin/killall

start()
{
	echo 1 > /proc/sys/net/br_igmp/igmp_snooping_switch
	/usr/sbin/mldsnoop &
}

stop()
{
	echo 0 > /proc/sys/net/br_igmp/igmp_snooping_switch
	$KILLALL mldsnoop
}
