# apv5sdk-v15/apps/sysstat-6.0.1/sysstat.cron.hourly
#!/bin/sh
# Run system activity accounting tool every 10 minutes
../../tools/ap71-rootdir//lib/sa/sa1 -d 600 6 &
